# templates/disable-inactive-freeipa-users.j2.service
[Unit]
Description=Disable inactive FreeIPA users.

[Service]
ExecStart=/usr/local/sbin/disable_inactive_freeipa_users.sh \
"{{ freeipa_days_before_inactive }} days ago"
# TODO: There are other sandboxing-related directives that could
# easily be added here.  See #75 for more details.
PrivateTmp=true
Type=oneshot
