# files/disable-inactive-freeipa-users.service
[Unit]
Description=Disable inactive FreeIPA users.

[Service]
ExecStart=/usr/local/sbin/disable_inactive_users.sh
PrivateTmp=true
Type=oneshot
