# files/disable-inactive-freeipa-users.timer
[Unit]
Description=Disable inactive FreeIPA users daily.

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
